# src/root/isync.cron
*/5 * * * * mbsync -c /mbsync.conf {NAME} >/dev/null 2>&1
