# src/root/etc/crontabs/root
*/5 * * * * mbsync -c /mbsyncrc.conf {NAME} >/dev/null 2>&1
